# deploy/systemd/news-agent.service
[Unit]
Description=AI News Agent (one-shot daily run)
After=network-online.target
Requires=ai-provider-service.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/ai-provider-service
EnvironmentFile=/var/www/ai-provider-service/.env
ExecStart=/var/www/ai-provider-service/venv/bin/python -m agents.news.runner
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
