# deploy/open-webui/litellm.container
# Quadlet: LiteLLM-Proxy. Übersetzt OpenAI-Wire (von Open WebUI) → natives Anthropic.
#
# Datei nach /etc/containers/systemd/ kopieren, dann:
#   systemctl daemon-reload && systemctl start litellm.service

[Unit]
Description=LiteLLM proxy for Open WebUI
After=network-online.target postgres.service
Wants=network-online.target
Requires=postgres.service

[Container]
Image=ghcr.io/berriai/litellm:main-stable
ContainerName=litellm
Network=open-webui.network
EnvironmentFile=/opt/open-webui/.env
Volume=/opt/open-webui/litellm-config.yaml:/app/config.yaml:ro
# Single-User-Setup: 1 Worker reicht. Spart ~1.2 GB RAM gegenüber 2 Workern,
# weil jeder Worker die volle LiteLLM-Lib im Speicher hält.
Exec=--config /app/config.yaml --port 4000 --num_workers 1
# Hard limit — OOM-Killer trifft nur den Container, nicht den Host.
Memory=1500m
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=default.target
