# deploy/open-webui/open-webui.container
# Quadlet: Open WebUI. An 127.0.0.1:3000 gebunden, Apache vorgeschaltet.
#
# Datei nach /etc/containers/systemd/ kopieren, dann:
#   systemctl daemon-reload && systemctl start open-webui.service

[Unit]
Description=Open WebUI
After=network-online.target litellm.service
Wants=network-online.target
Requires=litellm.service

[Container]
Image=ghcr.io/open-webui/open-webui:main
ContainerName=open-webui
Network=open-webui.network
PublishPort=127.0.0.1:3000:8080
# Persistente Daten (Chat-Verlauf, RAG-DB, Settings) als Host-Bind-Mount —
# einfacher zu backuppen als ein named volume.
Volume=/opt/open-webui/data:/app/backend/data:Z
EnvironmentFile=/opt/open-webui/.env
# WICHTIG: DATABASE_URL aus dem geteilten .env überschreiben, sonst migriert
# Open WebUI seine ganze Schema-Hierarchie zu Postgres (vorgesehen für LiteLLM)
# und verliert den Anschluss an die SQLite-Daten in /opt/open-webui/data/.
Environment=DATABASE_URL=sqlite:////app/backend/data/webui.db
Environment=WEBUI_AUTH=true
Environment=ENABLE_SIGNUP=false
Environment=DEFAULT_USER_ROLE=pending
# Anthropic Claude über LiteLLM-Proxy (OpenAI-kompatibel)
Environment=OPENAI_API_BASE_URL=http://litellm:4000/v1
Environment=ENABLE_OPENAI_API=true
# Ollama vom Mac via autossh-Reverse-Tunnel: VPS:11434 → Mac:11434.
# host.containers.internal löst aus dem Container den VPS-Host auf.
Environment=ENABLE_OLLAMA_API=true
Environment=OLLAMA_BASE_URL=http://host.containers.internal:11434
# Default-Modell für neue Chats — primär lokales Ollama, Claude als Option.
Environment=DEFAULT_MODELS=hf.co/Jiunsong/supergemma4-26b-uncensored-gguf-v2:Q4_K_M
# Hard limit — RAG-Embeddings können kurz spiken, deshalb großzügig.
Memory=3g
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=default.target
