# deploy/ai-provider-summary.timer
[Unit]
Description=ai-provider-service nightly summary job timer

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=ai-provider-summary.service

[Install]
WantedBy=timers.target
