# deploy/ai-provider-vault-render.timer
[Unit]
Description=ai-provider-service vault self-heal timer

[Timer]
OnUnitActiveSec=10min
OnBootSec=2min
Unit=ai-provider-vault-render.service

[Install]
WantedBy=timers.target
