# deploy/ai-provider.service
[Unit]
Description=AI Provider Service (Multi-Provider Chat Gateway)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/ai-provider-service
# bash-Wrapper notwendig wegen SELinux-Behandlung von venv-Symlinks (gleiches
# Pattern wie bewerbungen.service). .env wird von python-dotenv geladen.
ExecStart=/bin/bash -c 'exec /var/www/ai-provider-service/venv/bin/gunicorn --workers 2 --worker-class sync --bind 127.0.0.1:8767 --timeout 120 --access-logfile /var/log/ai-provider/access.log --error-logfile /var/log/ai-provider/error.log "app:create_app()"'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
