# deploy/ai-provider.container
# SPDX-License-Identifier: AGPL-3.0-or-later
# © 2026 Harald Weiss
#
# NOTE: This Quadlet is OBSOLETE for current VPS deployment.
# The service now runs via /etc/systemd/system/ai-provider.service
# (see ai-provider.service in this directory) because podman 4→5
# broke the podman1 bridge IP (10.88.0.1 → 10.89.0.1) and :Z volume
# labels cause MCS mismatch. The systemd unit uses --security-opt
# label=disable instead.
#
# Keep this file as reference for future Quadlet-based deploys once
# podman/podman-compat is fixed.

[Unit]
Description=AI Provider Service (Multi-Provider Chat Gateway)
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/ai-provider:latest
ContainerName=ai-provider

PublishPort=127.0.0.1:8767:8767

# SELinux label disabled — MCS mismatch between podman 5 bridge and :Z/:z
SecurityLabelDisable=true

# Persistent SQLite DB + pricing overrides outside the image.
Volume=/opt/ai-provider-data:/app/data
Volume=/opt/ai-provider-data/pricing_overrides.json:/app/pricing_overrides.json

HealthCmd=curl -fsS http://127.0.0.1:8767/health
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s

EnvironmentFile=/etc/ai-provider/ai-provider.env

Environment=DATABASE_URL=sqlite:////app/data/storage.db
Environment=TZ=Europe/Berlin
Environment=HOST=0.0.0.0
Environment=PORT=8767
Environment=GATE_ENABLED=true
Environment=UNGATED_PROVIDERS=ollama
Environment=HEALTH_CHECK_INTERVAL_SEC=30
Environment=QUEUE_DRAIN_INTERVAL_SEC=60
Environment=QUEUE_TTL_HOURS=24

Memory=512m
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
