# deploy/ai-provider.service
[Unit]
Description=AI Provider Service (Multi-Provider Chat Gateway)
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=/usr/bin/podman rm -f ai-provider
ExecStartPre=-/usr/sbin/fuser -k 8767/tcp
ExecStart=/usr/bin/podman run --name ai-provider --security-opt label=disable \
  -v /opt/ai-provider-data/storage.db:/app/data/storage.db \
  -v /var/lib/ai-provider-service/vault:/var/lib/ai-provider-service/vault \
  --publish 127.0.0.1:8767:8767 \
  --env-file /etc/ai-provider/ai-provider.env \
  localhost/ai-provider:latest
ExecStop=/usr/bin/podman stop -t 3 ai-provider
ExecStopPost=-/usr/sbin/fuser -k 8767/tcp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
